# MonitorizacionIOT_1.0_all/etc/systemd/system/downloadData.service
[Service]

ExecStart=/usr/bin/DescargaData/Download_data.sh
