# MonitorizacionIOT_1.0_all/etc/systemd/system/Monitorizacion_IOT.service
[Service]

ExecStart=/usr/bin/Monitorizacion_IOT/Monitorizacion_IOT.sh